# backup_daemon.service
[Unit]
Description=Backup Daemon

[Service]
ExecStart=/home/stepa/daemon/daemon
Restart=always

[Install]
WantedBy=default.target